# builder/os/mkosi.skeleton/etc/systemd/system/first-boot.service
[Unit]
Description=First Boot Configuration
After=multi-user.target
Wants=multi-user.target
ConditionPathExists=!/var/lib/first-boot-complete

[Service]
Type=oneshot
ExecStart=/opt/postinstall/first-boot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=default.target
