# generate-cloud-config/files/etc/profile.d/custom-settings.sh
#executable
#!/usr/bin/env bash
# System-wide custom settings
export HISTSIZE=10000
export HISTFILESIZE=20000